# deploy/.service
[Unit]
Description=Gunicorn instance to serve  {{app}}
After=network.target

[Service]
User={{ansible_ssh_user}}
Group=www-data
WorkingDirectory=/home/{{ansible_ssh_user}}/{{app}}
Environment="PATH=/home/{{ansible_ssh_user}}/{{app}}/venv/bin"
ExecStart=/home/{{ansible_ssh_user}}/{{app}}/venv/bin/gunicorn --workers 3 --bind unix:{{app}}.sock -m 007 run:app

[Install]
WantedBy=multi-user.target
